# supervisord.conf
[supervisord]
nodaemon=true

[program:mongo]
command=/usr/bin/mongod

[program:node]
command=/usr/bin/node /srv/server/server.js 